# examples/systemd/rsk-client.service
[Unit]
Description=RSK Client - Reverse SOCKS5 Proxy Client
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/tbxark/rsk

[Service]
Type=simple
User=rsk
Group=rsk
WorkingDirectory=/opt/rsk

# Client configuration
Environment="RSK_SERVER=your-server.com:7000"
Environment="RSK_TOKEN=CHANGE_ME_TO_SECURE_TOKEN"
Environment="RSK_PORTS=20001"
Environment="RSK_NAME=exit-node-1"

ExecStart=/opt/rsk/rsk-client \
  --server ${RSK_SERVER} \
  --token ${RSK_TOKEN} \
  --ports ${RSK_PORTS} \
  --name ${RSK_NAME} \
  --dial-timeout 15s \
  --allow-private-networks

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/rsk

# Restart policy (client auto-reconnects, but restart on crashes)
Restart=on-failure
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rsk-client

[Install]
WantedBy=multi-user.target
